# etc/profile.d/sencha.sh
#!/bin/bash
export PATH=/usr/bin/Sencha/Cmd:$PATH
